# Library/LaunchAgents/com.fredsh2k.notify-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fredsh2k.notify-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/fsherman/Code/Personal/dotfiles/bin/notify-daily</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/fsherman/.local/share/notify-daily/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/fsherman/.local/share/notify-daily/stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
