# Library/LaunchAgents/com.fredsh2k.briefing-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fredsh2k.briefing-morning</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/fsherman/Code/Personal/dotfiles/bin/briefing-morning</string>
  </array>

  <!-- Sun-Thu (Weekday: 0=Sun, 1=Mon, ..., 4=Thu) at 09:00 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>0</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/fsherman</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/fsherman/.local/share/briefing-morning/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/fsherman/.local/share/briefing-morning/stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
